# docker/etc/supervisor/conf.d/app.conf
[program:app]
command = chat_take_aiohttp %(ENV_COMMAND_NAME)s
directory = /app
redirect_stderr = true
stdout_logfile = /var/log/app.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 0
